# roles/tftpd/templates/tftpd.service.j2
[Unit]
Description=Dnsmasq TFTP Daemon

[Service]
ExecStart=/usr/sbin/dnsmasq -k --log-facility=- --port=0 --enable-tftp --tftp-root={{tftpd_dir}}
Restart=always
