# roles/tftpd/templates/tftpd.service.j2
[Unit]
Description=TFTP Daemon

[Service]
ExecStart=/usr/bin/in.tftpd -s {{ tftpd_dir }}
StandardInput=socket
StandardOutput=inherit
StandardError=journal
