# roles/tftpd/templates/tftpd-announce.service.j2
[Unit]
Description=Announce tftpd
BindsTo=tftpd.service
After=tftpd.service

[Service]
ExecStart=/bin/sh -c "while true; do /usr/bin/etcdctl set /lochness/services/tftpd/{{machine_id}} '{ \"target\": \"{{machine_id}}.nodes.\", \"port\": 67 }' --ttl 60 > /dev/null; sleep 30; done"
ExecStop=/usr/bin/etcdctl rm /lochness/services/tftpd/{{machine_id}}
User=nobody
RestartSec=10s
Restart=always
