# roles/dhcpd/templates/dhcpd.service.j2
[Unit]
Description=IPv4 DHCP server
Wants=network.target
After=network.target
Wants=dhcpd-announce.service

[Service]
ExecStart=/usr/bin/dhcpd -4 -d -q -f -cf /etc/dhcpd/dhcpd.conf
StandardOutput=inherit
StandardError=journal
RestartSec=5s
Restart=always

[Install]
WantedBy=multi-user.target
