# roles/cguestd/files/cguestd.service
[Unit]
Description=Lochness Guest Management
Wants=network.target etcd.service
After=network.target etcd.service

[Service]
ExecStart=/usr/sbin/cguestd
StandardOutput=journal
StandardError=journal
Restart=always

[Install]
WantedBy=multi-user.target
